# deploy/seiche-publish.service
[Unit]
Description=Seiche hourly forecast publish
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=buoycast
WorkingDirectory=/opt/seiche
Environment=HOME=/opt/seiche
ExecStart=/opt/seiche/venv/bin/python publish.py
TimeoutStartSec=900
Nice=10
